# etc/cron.d/zabbix
MAILTO=root

*/1 * * * * root [ -x /etc/zabbix/bin/php-fpm-status-collector.sh ] && /etc/zabbix/bin/php-fpm-status-collector.sh | jq -c > /tmp/php-fpm-status-collector.json
